# xcp-orchestrator/systemd/xcp-orchestrator.timer
[Unit]
Description=Trigger XCP-ng Agentic Orchestration Daily at 5 AM

[Timer]
OnCalendar=*-*-* 05:00:00
Persistent=true
Unit=xcp-orchestrator.service

[Install]
WantedBy=timers.target
