# xcp-orchestrator/systemd/xcp-orchestrator.service
[Unit]
Description=XCP-ng Agentic Orchestration Automation Run
After=network-online.target ollama.service

[Service]
Type=oneshot
User=root
# Securely injects GitHub API Token using systemd sandbox file allocation
LoadCredential=GITHUB_TOKEN:/etc/xcp-hl-credentials/github_token
ExecStart=/usr/local/bin/xcp-orchestrator

[Install]
WantedBy=multi-user.target
